# deploy/aws/neurotrust-ms.service
[Unit]
Description=NeuroTrust-MS FastAPI backend
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/neurotrust-ms/backend
EnvironmentFile=/etc/neurotrust-ms.env
Environment=PYTHONPATH=/opt/neurotrust-ms/backend
ExecStart=/opt/neurotrust-ms/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
